# .runit/sv/12to11/run
#!/bin/sh

[ -r ./conf ] && . ./conf

export $(dbus-launch)
exec 12to11 2>&1

