# .runit/sv/pulseaudio/run
#!/bin/sh

[ -r ./conf ] && . ./conf

exec pulseaudio ${OPTS:-} 2>&1
